# configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/run/supervisord.pid
user=root

[unix_http_server]
file=/run/supervisor.sock
chmod=0700
username=root
password=root

[supervisorctl]
serverurl=unix:///run/supervisor.sock
username=root
password=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:sshd]
command=/usr/bin/sshd -D -e
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.err

[program:ttyd]
command=/usr/bin/ttyd -W -t 'unicodeVersion=11' -i 127.0.0.1 /tmux.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/ttyd.log
stderr_logfile=/var/log/ttyd.err

[program:webui]
command=/usr/bin/gunicorn --chdir /root/webui --bind 127.0.0.1:8983 --workers 1 --threads 8 --timeout 0 --graceful-timeout 30 --access-logfile - --error-logfile - app:create_wsgi_app()
environment=CODEX_HOME="/data/codex",AUDITOR_WEBUI_STATIC="/root/webui/static"
autostart=true
autorestart=true
stdout_logfile=/var/log/webui.log
stderr_logfile=/var/log/webui.err

[program:nginx]
command=/usr/bin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.err

[program:filebrowser]
command=/usr/bin/filebrowser -d /data/filebrowser.db
autostart=true
autorestart=true
stdout_logfile=/var/log/filebrowser.log
stderr_logfile=/var/log/filebrowser.err
